# distro/debian-live/config/includes.chroot/etc/systemd/system/seatd.service
[Unit]
Description=Seat management daemon
Documentation=man:seatd(1)
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/bin/seatd -g seat
Restart=on-failure

[Install]
WantedBy=multi-user.target
